# src/webapp/bin/webapp
#!/bin/bash
# /etc/init.d/webapp

### BEGIN INIT INFO
# Provides:          webapp
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Webapp for the charcuterie controller
# Description:       A mobile app to view and control the charcuterie controller
### END INIT INFO


case "$1" in 
    start)
        echo "Starting webapp"
        cd /home/pi/webapp/app/bin
	./webapp.d &
        ;;
    stop)
        echo "Stopping webapp"
        killall -g -w webapp.d
        ;;
    *)
        echo "Usage: /etc/init.d/webapp start|stop"
        exit 1
        ;;
esac

exit 0
